# .config/systemd/user/obsidian-sync.timer
[Unit]
Description=Run Obsidian sync periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target