# .config/systemd/user/obsidian-sync.timer
[Unit]
Description=Run Obsidian sync periodically

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target