# .config/systemd/user/obsidian-sync.service
[Unit]
Description=Sync Obsidian vault with Google Drive (rclone bisync)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/sync-obsidian.sh